# upstart/lxc-autostart.conf
# lxc-autostart - starts all containers that want to autostart on boot

description	"start containers"

start on runlevel [2345]

task

script
    exec lxc autostart
end script
