# upstart/lxc-stopall.conf
# lxc-stopall - stops all containers when shutting down

description	"stop containers"

start on runlevel [016]

task

script
    exec lxc stopall
end script
